# lib/supervisord/moodsy_www.conf
[program:moodsy_www]
directory=/var/www/python/moodsy.me/moodsy_www
command=python3 main.py
stopsignal=KILL
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/moodsy_www
environment=
    MOODSY_WWW_STATIC_PATH="/var/www/python/moodsy.me/static",
    MOODSY_WWW_TEMPLATE_PATH="/var/www/python/moodsy.me/templates",
    MOODSY_WWW_SOCKET="/var/www/sockets/moodsy_www.sock",
    PATH="env/bin"
